# pisound-btn/debian/DEBIAN/postinst
systemctl daemon-reload
if [ ! -d "/usr/local/pisound" ]; then
	echo "Cloning pisound repository from https://github.com/BlokasLabs/pisound..."
	git clone https://github.com/BlokasLabs/pisound.git /usr/local/pisound
else
	echo "Updating pisound repository with latest stuff in https://github.com/BlokasLabs/pisound..."
	cd /usr/local/pisound && git -c user.name="apt-get" -c user.email="apt@get" stash && git pull
fi

if [ -e "/usr/local/bin/pisound-btn" ]; then
	echo "Deleting locally built pisound-btn, so that the packaged version is preferred."
	rm -f "/usr/local/bin/pisound-btn"
fi

if [ -e "/etc/systemd/system/pisound-btn.service" ]; then
	echo "Deleting /etc/systemd/system/pisound-btn.service from 1.04"
	systemctl stop pisound-btn
	systemctl disable pisound-btn
	rm -f "/etc/systemd/system/pisound-btn.service"
fi

# Automatically added by dh_installinit
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ]; then
	systemctl enable pisound-btn
	systemctl start pisound-btn
fi
# End automatically added section

# Check if we're being upgraded on MODEP image.
if [ -e "/usr/lib/systemd/system/pisound-hotspot.service" ]; then
	sed -i 's/Pisound/MODEP/g' /etc/hostapd/hostapd.conf
	if $(systemctl is-enabled -q pisound-hotspot); then
		systemctl disable pisound-hotspot
		systemctl stop pisound-hotspot
		systemctl enable wifi-hotspot
		systemctl start wifi-hotspot
	fi
	rm -f /usr/lib/systemd/system/pisound-hotspot.service
	systemctl daemon-reload
else
	systemctl daemon-reload
	if $(systemctl is-enabled -q wifi-hotspot); then
		systemctl restart wifi-hotspot
	fi
fi

if [ -e /etc/pisound.conf ]; then
	if [ ! -L /etc/pisound.conf ]; then
		mv /etc/pisound.conf /usr/local/etc/pisound.conf.bak.`date +"%F.%T"`
	else
		unlink /etc/pisound.conf
	fi
	ln -s /usr/local/etc/pisound.conf /etc/pisound.conf
fi
